# docker/supervisor.conf.d/janitoo_nut.conf
[program:jnt_nut]
user=root
autostart=true
autorestart=true
command=jnt_nut -c /opt/janitoo/etc/janitoo_nut.conf start
