# docker/supervisor.conf.d/janitoo_nut.conf
[program:jnt_nut]
user=root
autostart=true
autorestart=true
command=jnt_nut -c /etc/janitoo/janitoo_nut.conf front
